# templates/default/golang-agent-init.erb
#!/bin/bash
#
# File autogenerated by chef. All changes will be lost
#
### BEGIN INIT INFO
# Provides: <%= @agent_name %>
# Required-Start: $network $remote_fs
# Required-Stop: $network $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Start the GOCD Agent
### END INIT INFO
NAME="<%= @agent_name %>"
PIDDIRECTORY="/var/run/gocd-golang-agent"
PIDFILE="/var/run/gocd-golang-agent/<%= @agent_name %>.pid"
DAEMON="/usr/bin/gocd-golang-agent"

. /etc/default/<%= @agent_name %>
<% if @autoregister_file %>
. <%= @autoregister_file %>
<% end %>

case "$1" in
  start)
        echo -n "Starting daemon: "$NAME
  mkdir -p $PIDDIRECTORY
	/sbin/start-stop-daemon --start -c go:go -m --quiet --pidfile $PIDFILE --background --exec $DAEMON
        echo "."
	;;
  stop)
        echo -n "Stopping daemon: "$NAME
	/sbin/start-stop-daemon --stop -c go:go --quiet --oknodo --pidfile $PIDFILE
        echo "."
	;;
  restart)
        echo -n "Restarting daemon: "$NAME
  mkdir -p $PIDDIRECTORY
	/sbin/start-stop-daemon --stop -c go:go --quiet --oknodo --retry 30 --pidfile $PIDFILE
	/sbin/start-stop-daemon --start -c go:go -m --quiet --pidfile $PIDFILE --background --exec $DAEMON -- $DAEMON_OPTS
	echo "."
	;;

  *)
	echo "Usage: "$1" {start|stop|restart}"
	exit 1
esac
exit 0
